# systemd/blood_moon.service
[Unit]
Description=Run bm.py with uv

[Service]
Type=simple
WorkingDirectory=/path/to/your/project
ExecStartPre=/usr/bin/uv sync
ExecStart=/usr/bin/uv run blood_moon